# cookbooks/redistest/templates/default/redis-sentinel.upstart.erb
# this file is generated by Chef.  Manual changes will be overwritten

description "redis sentinel"

start on started redis
stop on shutdown

exec /usr/local/bin/redis-server /etc/redis-sentinel.conf --sentinel
